# rootfs/etc/services.d/tailscaled/run
#!/usr/bin/execlineb -P
tailscaled --state /state/tailscale.state