# rootfs/etc/services.d/tailscale-web/run
#!/bin/sh
sleep 10
exec tailscale web -listen 0.0.0.0:8088